# iatis-api.service
[Unit]
Description=IATIS FastAPI Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# TODO(ops): run as a dedicated non-root user instead of root. This was not
# switched automatically because it requires VPS-side setup first:
#   useradd -r -s /usr/sbin/nologin iatis
#   chown -R iatis:iatis /root/IATIS
# then set User=iatis / Group=iatis here. Until that's done, the hardening
# directives below still meaningfully reduce blast radius even under root.
User=root
WorkingDirectory=/root/IATIS
Environment=PATH=/root/IATIS/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
# python3 -m uvicorn (not the venv/bin/uvicorn console script): survives
# venv rebuilds where the entry-point script is missing — systemd reports
# those as status=203/EXEC before Python even starts.
ExecStart=/root/IATIS/venv/bin/python3 -m uvicorn execution.api_server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Sandboxing — the API server accepts untrusted network input, so a bug or
# RCE here should not get free rein over the filesystem.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
# read-only (not true): the code and venv live under /root/IATIS, and
# ProtectHome=true hides /root entirely from the service namespace — the
# ExecStart binary itself becomes unreachable and systemd dies with
# status=203/EXEC before Python starts (observed on the VPS, 2026-07-05).
# read-only still blocks writes anywhere in /root except the
# ReadWritePaths=/root/IATIS/storage carve-out below.
ProtectHome=read-only
ReadWritePaths=/root/IATIS/storage
MemoryMax=1G
TasksMax=256

[Install]
WantedBy=multi-user.target
